# lib/profile.d/m7.sh
#!/bin/bash

# NOTE: You need to copy this file to the '/etc/profile.d' directory on your
# server for this platform to work when being run from cron or with init.d
# scripts. The location might change from server to server, so make sure you
# update this file accordingly.
export M7_ROOT="/opt/vpls/m7"